# hregionserver/supervisord.conf
[supervisord]
nodaemon=true

[program:hbase-regionserver]
command=service hbase-regionserver start
