# zk/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-hdfs-namenode]
command=/usr/bin/service zookeeper-server start
