# zk/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper-server]
command=/usr/bin/service zookeeper-server start
